# files/etc/crontabs/root
30 3 * * 1 sleep 70 && touch /etc/banner && reboot

45 18 * * * ash /etc/change2oraclennr.sh
15 1 * * * ash /etc/change2oraclehy.sh