# files/etc/crontabs/root
30 3 * * 1 sleep 70 && touch /etc/banner && reboot
